# Linux/etc/profile.d/btrfs.sh
export GRUB_BTRFS_LIMIT_TOTAL_ENTRIES=150
